# files/services/dnstap/run
#!/usr/bin/env sh
export GOPATH
exec su-exec unbound ${GOPATH}/bin/dnstap -u ${GOPATH}/var/run/dnstap.sock
