# files/services/dnstap/run
#!/usr/bin/env sh
export GOPATH
exec gosu unbound:unbound bash -c "${GOPATH}/bin/dnstap -u ${GOPATH}/var/run/dnstap.sock"
